# meta-facebook/meta-tiogapass/recipes-phosphor/flash/phosphor-software-manager/obmc-flash-host-bios@.service
[Unit]
Description=Flash Host Bios image %I to Host

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/home/root/bios-update.sh /tmp/image/%i
