# recipes-phosphor/sensors/dbus-sensors/powerOffSlot.service
[Unit]
Description=Power off the slots when fan fails.

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/sbin/powerOffSlot.sh
