# recipes-phosphor/datetime/phosphor-time-manager/bmc-set-time.service
[Unit]
Description=get time from host and set bmc time

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/sbin/SetTimeBmc.sh
